# etc/xdg/autostart/thymetrack.desktop
[Desktop Entry]
Type=Application
Name=Thyme Track
Exec=thymetrack
StartupNotify=false
Terminal=false
